# deploy/vps/cyberos.service
[Unit]
Description=CyberOS Docker Compose stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/cyberos/repo
RemainAfterExit=yes
ExecStart=/usr/bin/docker compose --env-file deploy/vps/.env -f deploy/vps/compose.prod.yml up -d --build
ExecStop=/usr/bin/docker compose --env-file deploy/vps/.env -f deploy/vps/compose.prod.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
